# compose/slave.conf
[supervisord]
nodaemon=true

[program:spark-slave]
command=/opt/spark/sbin/start-slave.sh -c 1 spark://master:7077
